# apps/market_radar/scripts/com.eniskorkut.bist-radar.telegram.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.eniskorkut.bist-radar.telegram</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/eniskorkut/Documents/New project 2/bist-research/apps/market_radar/scripts/run_telegram_radar_launchd.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>15</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>18</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/eniskorkut/Documents/New project 2/bist-research</string>

  <key>StandardOutPath</key>
  <string>/Users/eniskorkut/Documents/New project 2/bist-research/logs/telegram_radar.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/eniskorkut/Documents/New project 2/bist-research/logs/telegram_radar.launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
